# chezmoi/Library/LaunchAgents/com.jiz4oh.lemonade.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.jiz4oh.lemonade</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>/Users/jiz4oh/bin/lemonade-daemon</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/jiz4oh/Library/Logs/lemonade/launchd.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/jiz4oh/Library/Logs/lemonade/launchd.out.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/jiz4oh</string>
</dict>
</plist>
