# chezmoi/Library/LaunchAgents/com.minis.tg.listen.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.minis.tg.listen</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>{{ .chezmoi.homeDir }}/bin/tg-listen-daemon</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>{{ .chezmoi.homeDir }}/Library/Logs/tg-cli/launchd.err.log</string>
	<key>StandardOutPath</key>
	<string>{{ .chezmoi.homeDir }}/Library/Logs/tg-cli/launchd.out.log</string>
	<key>WorkingDirectory</key>
	<string>{{ .chezmoi.homeDir }}</string>
</dict>
</plist>

